# deploy/sql-mastery.service
# systemd unit for SQL Mastery Path.
# Install:
#   sudo cp deploy/sql-mastery.service /etc/systemd/system/sql-mastery.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now sql-mastery

[Unit]
Description=SQL Mastery Path
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=scott
Group=scott
WorkingDirectory=/home/scott/sql-mastery
EnvironmentFile=/home/scott/sql-mastery/.env
ExecStart=/usr/bin/node /home/scott/sql-mastery/dist/server.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
